# main/.config/systemd/user/udiskie.service
[Unit]
Description=Automatically mount newly inserted USB drives

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
